# macos/peon/com.peon-ping.relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template. install.sh substitutes __PEON_BIN__ and __HOME__ at install
  time so the brew prefix (Intel vs Apple Silicon) resolves correctly.
  Do not symlink this file — it is rendered into ~/Library/LaunchAgents/.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.peon-ping.relay</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PEON_BIN__</string>
        <string>relay</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/peon-relay.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/peon-relay.err.log</string>
</dict>
</plist>
